# external/board/mango-rpi4-v1/overlay/etc/cron.d/exchange-watchdog
# Reboot if no connection to outside world 
0 * * * * /usr/local/bin/exchange-watchdog
