# external/board/mango-rpi4-v1/overlay/etc/cron.d/flag-watchdog
# Take action if various flag files are set
* * * * * /usr/local/bin/flag-watchdog
