# external/board/mango-rpi4-v1/overlay/etc/cron.d/mgetty-watchdog
# Kill stale mgetty processes
* * * * * /usr/local/bin/mgetty-watchdog
