# external/board/mango-rpi4-v1/overlay/etc/cron.d/tincan
# Manage tincan network 
*/1 * * * * /usr/local/bin/tincanctl -i serialnum cron
