# external/board/mango-rpi4-v1/overlay/etc/cron.d/update
# Check for updates 
*/1 * * * * /usr/local/bin/check_updates
