# external/board/mango-rpi4-v1/overlay/etc/profile.d/mango.sh
export PATH=$PATH:/opt/transport/groups/mango/support/bin
export PYTHONPATH=/opt/transport/groups/mango/support/lib:/opt/transport/lib/python2.7/site-packages
